# root/etc/s6-overlay/s6-rc.d/init-mod-dev_1c-install/up
/etc/s6-overlay/s6-rc.d/init-mod-dev_1c-install/run